# droplet.service
[Unit]
Description=droplet.service
After=syslog.target network-online.target

[Service]
Environment=GOTRACEBACK=crash
LimitCORE=infinity
LimitNOFILE=65536
WorkingDirectory=/tmp/
ExecStart=/usr/sbin/droplet
Restart=on-failure
RestartSec=10
MemoryLimit=20G

[Install]
WantedBy=multi-user.target
